# deploy/systemd/beemax.target
[Unit]
Description=BeeMax Agent profiles
Wants=network-online.target
After=network-online.target

[Install]
WantedBy=multi-user.target
