# deploy/systemd/beemax@.service
[Unit]
Description=BeeMax Agent profile %i
After=network-online.target
Wants=network-online.target
PartOf=beemax.target

[Service]
Type=simple
User=beemax
Group=beemax
WorkingDirectory=/opt/beemax
Environment=NODE_ENV=production
Environment=BEEMAX_PROFILE=%i
EnvironmentFile=-/etc/beemax/%i.env
ExecStart=/usr/bin/env npm run gateway -- --profile %i
Restart=on-failure
RestartSec=5s
TimeoutStopSec=60s
KillSignal=SIGTERM
NoNewPrivileges=true
PrivateTmp=true
UMask=0077
StandardOutput=journal
StandardError=journal
SyslogIdentifier=beemax-%i

[Install]
WantedBy=multi-user.target
