# deploy/systemd/beemax@.service
[Unit]
Description=BeeMax Agent profile %i
After=network-online.target
Wants=network-online.target
PartOf=beemax.target

[Service]
Type=simple
User=beemax
Group=beemax
WorkingDirectory=/opt/beemax
EnvironmentFile=-/var/lib/beemax/profiles/%i/.env
EnvironmentFile=-/opt/beemax/config/profiles/%i.env
EnvironmentFile=-/etc/beemax/%i.env
Environment=NODE_ENV=production
Environment=BEEMAX_PROFILE=%i
ExecStart=/usr/bin/env npm run gateway -- --profile %i --home /var/lib/beemax --root /opt/beemax
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=60
StartLimitBurst=5
TimeoutStopSec=60s
KillSignal=SIGTERM
NoNewPrivileges=true
PrivateTmp=true
UMask=0077
StandardOutput=journal
StandardError=journal
LogRateLimitIntervalSec=30s
LogRateLimitBurst=200
SyslogIdentifier=beemax-%i

[Install]
WantedBy=multi-user.target
